# supervisord.conf
[supervisord]
nodaemon=true

[program:node]
command=node /usr/share/nginx/html/index.js
autorestart=true

[program:nginx]
command=nginx -g 'daemon off;'
autorestart=true


